# docker/php/supervisor.ini
[program:runitonce]
process_name=%(program_name)s
command=php /app/artisan queue
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/runitonce.log
